# dist/openrc-init-script
#!/sbin/openrc-run

command="/usr/bin/uwsgi"

command_args="-s /var/lib/qdbs/uwsgi.sock --manage-script-name --mount /=pyqdbs:app"

command_dir="/var/lib/qdbs/pyqdbs/"

command_user="qdbs"

command_background="yes"
pidfile="/run/${SVCNAME}.pid"
start_stop_daemon_args="--chdir ${command_dir} -w 400"

depend()
{
        use     dns net time
}


start_post()
{
	chmod g+rwx /var/lib/qdbs/uwsgi.sock
}